# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# Fix Luci Don't Remove!!
umount /etc/resolv.conf
rm /etc/resolv.conf
ln -s /tmp/resolv.conf.auto /etc/resolv.conf

if [ -n "$root_password" ]; then
  (echo "$root_password"; sleep 1; echo "$root_password") | passwd > /dev/null
fi

exit 0
